# ar 5315u/custom-files/etc/init.d/telnetd
#!/bin/sh /etc/rc.common
# [2024-06-24] elias
START=90

start() {
        telnetd -l /usr/bin/vtysh
}

stop() {
        killall telnetd 2>/dev/null
}


